# recipes-bsp/bt-start/files/mcm-imx8m-mini/bt-start.service
[Unit]
Description=CompuLab BLueThooth Start Script

[Service]
Type=forking
ExecStart=/usr/sbin/bt-start start
ExecStop=/usr/sbin/bt-start stop

[Install]
WantedBy=multi-user.target
